# net-im/pyaim-t/files/pyaim-t-0.7c.initd
#!/sbin/runscript

depend() {
	need net
	need jabber-server
}

start() {
	ebegin "Starting AIM Jabber Transport"
	cd /var/run/jabber
	/usr/lib/python2.4/site-packages/pyaim-t/PyAIMt.py -c /etc/jabber/pyaim-t.xml -b -p /var/run/jabber/pyaim-t.pid
	eend $?
}

stop() {
	ebegin "Stopping AIM Jabber Transport"
	start-stop-daemon --stop --quiet --pidfile /var/run/jabber/pyaim-t.pid
	eend $?
}
