# net-misc/pglinux/files/peerguardianrc
#!/sbin/runscript
# Created by iainel 2005.
# $Header:  $

OPTS="-d -c /usr/share/peerguardian/PG.conf -l /var/log/PG.log -h -m"

start() {
	ebegin "Starting PeerGuardian Daemon"
	start-stop-daemon --start --quiet --exec /usr/bin/peerguardnf -- ${OPTS}
	eend $?
}

stop() {
	ebegin "Stopping PeerGuardian Daemon"
	start-stop-daemon --stop --quiet --exec /usr/bin/peerguardnf	
	eend $?
}

restart() {
	svc_stop; svc_start
}
